# .edge/edge-cleanup.conf
[program:edge_cleanup]
command=/opt/selenium/scripts/edge_cleanup.sh
autostart=true
autorestart=true
stderr_logfile=/dev/null
redirect_stderr=true
